# currency-bot.conf
[program:currencybot]
directory=/home/ubuntu/bots/currencybot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/bots/currencybot/logfile_err.log
stdout_logfile=/home/ubuntu/bots/currencybot/logfile.log